# packaging/systemd/sage-gloved.service.in
[Unit]
Description=Sage Glove agent-session supervisor
Documentation=https://github.com/sage-protocol/sage-protocol

[Service]
Type=simple
ExecStart=@GLOVED_BINARY@ --runtime-dir @GLOVE_RUNTIME_DIR@ --audit-key @GLOVE_AUDIT_KEY@ --journal @GLOVE_RECEIPT_JOURNAL@ --session-policy @GLOVE_SESSION_POLICY@ --session-store @GLOVE_SESSION_STORE@ --materialization-root @GLOVE_MATERIALIZATION_ROOT@ --library-bundle-root @GLOVE_LIBRARY_BUNDLE_ROOT@ --path-exposure-policy @GLOVE_PATH_EXPOSURE_POLICY@ --path-exposure-journal @GLOVE_PATH_EXPOSURE_JOURNAL@
Restart=on-failure
RestartSec=2s
TimeoutStopSec=15s
UMask=0077

[Install]
WantedBy=default.target
